# cloudflare_ip_updater/cf_ip_updater.timer
[Unit]
Description=cf_ip_updater	

[Timer]
OnUnitActiveSec=3600s
OnBootSec=3600s

[Install]
WantedBy=timers.target
